# puppet/modules/chassis/files/chassis-hosts.service
[Unit]
Description=Chassis subdomain hosts updater
StartLimitBurst=5

[Service]
TimeoutStartSec=0
TimeoutStopSec=600
Restart=on-failure
RestartSec=30s
ExecStart=/usr/bin/sudo avahi-publish-aliases

[Install]
WantedBy=vagrant.mount
