# puppet/modules/chassis/files/chassis-hosts.service
[Unit]
Description=Avahi/mDNS CNAME publisher
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mdns-publish-cname --ttl 20 vagrant.local wat.vagrant.local eep.vagrant.local cool.vagrant.local nice.vagrant.local
Restart=no
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
