# puppet/modules/chassis/files/chassis-hosts.service
[Unit]
Description=Chassis subdomain hosts updater

[Service]
TimeoutStartSec=0
TimeoutStopSec=600
Restart=on-failure
RestartSec=30s
ExecStart=/usr/bin/python /vagrant/puppet/chassis-hosts.py
